# rootfs/etc/services.d/00-cron/run
#!/usr/bin/execlineb -P

# Importa $CRON, sino devuelve 0
importas -D 0 cron CRON

# Verfica si $CRON esta deshabilitado y termina el script
if { s6-test $cron -eq 0 }

# Ejecuta cron
with-contenv
/usr/sbin/cron -f -l 0
